# claude-manager.service
[Unit]
Description=Claude Account Manager — unified session renewal + smart swap
After=network.target

[Service]
Type=simple
ExecStart=/home/baconnier/dev/CCAutoRenew/claude-manager.sh
Restart=on-failure
RestartSec=120
Environment=PATH=/home/baconnier/.local/bin:/home/baconnier/.nvm/versions/node/v20.20.0/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
